# root/etc/s6-overlay/s6-rc.d/svc-m3u4prox/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app || exit 1

# Ensure EPG directories exist and have correct permissions
/bin/mkdir -p /data/epg/tmp /data/epg-sites /data/vod-strm 2>/dev/null || true

# Ensure symbolic links exist
[ ! -L /epg ] && /bin/ln -sf /data/epg /epg
[ ! -L /epg-sites ] && /bin/ln -sf /data/epg-sites /epg-sites

# Create symlink to node_modules so config files can import dependencies
[ ! -L /data/epg-sites/node_modules ] && /bin/ln -sf /app/node_modules /data/epg-sites/node_modules

# Set proper ownership on key directories only (avoid recursing through 480k+ STRM files)
/bin/chown -R m3u4prox:m3u4prox /data/db /data/epg /data/epg-sites /data/config 2>/dev/null || true
/bin/chown m3u4prox:m3u4prox /data/vod-strm 2>/dev/null || true
/bin/chmod -R u+w /data/epg-sites 2>/dev/null || true

echo "[m3u4prox] starting on port ${PORT:-3005}"

exec s6-setuidgid m3u4prox node /app/server/index.js
