# root/etc/s6-overlay/s6-rc.d/svc-m3u4prox/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Ensure EPG directories exist and have correct permissions
mkdir -p /data/epg/tmp /data/epg-sites 2>/dev/null || true

# Ensure symbolic links exist
[ ! -L /epg ] && ln -sf /data/epg /epg
[ ! -L /epg-sites ] && ln -sf /data/epg-sites /epg-sites

# Create symlink to node_modules so config files can import dependencies
[ ! -L /data/epg-sites/node_modules ] && ln -sf /app/node_modules /data/epg-sites/node_modules

# Set proper ownership - do this EVERY time to fix any root-owned subdirectories
chown -R m3u4prox:m3u4prox /data 2>/dev/null || true
chmod -R u+w /data/epg-sites 2>/dev/null || true

echo "[m3u4prox] starting on port ${PORT:-3005}"

exec s6-setuidgid m3u4prox node /app/server/index.js
